# static/curl/rc.local
#!/bin/sh
ip addr add 10.64.64.64/32 dev lo
